# etc/systemd/system/container-feedbin.service
[Unit]
Description=Podman container-feedbin.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target container-caddy.service container-feedbin-extract.service container-feedbin-privacy-please.service container-feedbin-rake.service container-mc@feedbin.service container-redis.service
After=network-online.target container-caddy.service container-feedbin-extract.service container-feedbin-privacy-please.service container-feedbin-rake.service container-mc@feedbin.service container-redis.service
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run --cidfile=%t/%n.ctr-id --sdnotify=conmon --cgroups=no-conmon --rm -d --replace --name feedbin --net host --volume /var/run/postgresql:/var/run/postgresql:rw --env-file /data/feedbin/.env --env RAILS_SERVE_STATIC_FILES=true localhost:5000/feedbin:edge /bin/sh -c 'bundle exec rake assets:precompile && bundle exec foreman start && exec bundle exec rails server'
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm -f --ignore --cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=multi-user.target
